# config/systemd/claude-chat.service
[Unit]
Description=PIBULUS Claude Deck Gateway
After=network.target

[Service]
Type=simple
User=pibulus
Group=pibulus
WorkingDirectory=/home/pibulus/pibulus-os
Environment=PYTHONUNBUFFERED=1
Environment=CLAUDE_CHAT_HOST=172.17.0.1
Environment=CLAUDE_CHAT_PORT=9016
Environment=CLAUDE_CHAT_ALLOWED_ORIGINS=https://deck.quickcat.club
Environment=CLAUDE_CHAT_ENABLE_FULL=0
Environment=CLAUDE_CHAT_MIN_AVAILABLE_MB=650
Environment=CLAUDE_CHAT_MAX_LOAD_1=5.0
Environment=CLAUDE_CHAT_REQUIRE_SWAP=1
EnvironmentFile=-/home/pibulus/.config/claude-chat.env
ExecStart=/usr/bin/python3 /home/pibulus/pibulus-os/scripts/claude_chat_gateway.py
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=15

NoNewPrivileges=true
UMask=0077
PrivateTmp=true
ProtectSystem=full
CapabilityBoundingSet=
RestrictSUIDSGID=true
LockPersonality=true
MemoryAccounting=yes
MemoryMax=768M
CPUAccounting=yes
CPUQuota=75%
TasksMax=64
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
