# config/systemd/deck-doctor.timer
[Unit]
Description=Run PIBULUS Deck doctor status pass

[Timer]
OnBootSec=4min
OnUnitActiveSec=15min
RandomizedDelaySec=90s
Unit=deck-doctor.service

[Install]
WantedBy=timers.target
