# config/systemd/isitgoingtorain.service
[Unit]
Description=Is It Going To Rain Static App
After=network.target

[Service]
Type=simple
User=pibulus
WorkingDirectory=/home/pibulus/apps/isitgoingtorain
ExecStart=/usr/bin/python3 -m http.server 9007 --bind 127.0.0.1 --directory /home/pibulus/apps/isitgoingtorain
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
