# config/systemd/kpab-services.service
[Unit]
Description=KPAB.FM Microservices (Deck + Shoutbox)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
KillMode=control-group
ExecStart=/bin/bash /home/pibulus/pibulus-os/scripts/start-kpab-services.sh
Restart=on-failure
RestartSec=30
User=pibulus

[Install]
WantedBy=multi-user.target
