# config/systemd/pibulus-watchdog.timer
[Unit]
Description=Run PIBULUS service health watchdog

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
RandomizedDelaySec=90s
Unit=pibulus-watchdog.service

[Install]
WantedBy=timers.target
