# config/systemd/qrbuddy.service
[Unit]
Description=QRBuddy Fresh App
After=network.target

[Service]
Type=simple
User=pibulus
WorkingDirectory=/home/pibulus/apps/qrbuddy
ExecStart=/home/pibulus/.deno/bin/deno run -A main.ts
Environment=PORT=8004
Environment=APP_URL=https://qrbuddy.app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
