# config/systemd/talktype.service
[Unit]
Description=TalkType SvelteKit App
After=network.target

[Service]
Type=simple
User=pibulus
WorkingDirectory=/home/pibulus/apps/talktype
ExecStart=/usr/bin/node index.js
Environment=PORT=9002
Environment=HOST=0.0.0.0
Environment=ORIGIN=https://talktype.app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
